# switchable-graphics/DEBIAN/postinst
# Install init script
sudo update-rc.d insert_radeon defaults 80

# Blacklist radeon driver
sed '$ a\
\
blacklist radeon' /etc/modprobe.d/blacklist.conf > blacklist.conf

sudo mv blacklist.conf /etc/modprobe.d/

